# etc/crontab
0-30 06 * * * flock -n /tmp/rss -c rss.sh/cron.sh
